# dot_config/systemd/user/joplin-to-git-backup.timer
[Unit]
Description=Joplin to Git backup timer

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
